# voice-wake.service
[Unit]
Description=Voice Wake Service for OpenClaw
After=sound.target network.target

[Service]
Type=simple
WorkingDirectory=/home/kaushik/.openclaw/workspace/claw-voice-wake/scripts
ExecStart=/home/kaushik/.openclaw/workspace/claw-voice-wake/scripts/voice_wake_wrapper.sh
Restart=always
RestartSec=5
Environment=HOME=/home/kaushik

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
